# system-integration/systemd/ss-susi-server.service.in
[Unit]
Description=Starting SUSI Server for SUSI Linux

[Service]
User=@SUSI_SERVER_USER@
Type=forking
PIDFile=@SUSIDIR@/susi_server/data/susi.pid
Environment=SS_LOG_TO_STDOUT=1
ExecStart=@SUSIDIR@/bin/susi-server start

[Install]
WantedBy=multi-user.target
