# systemd/ss-susi-server@.service.in
[Unit]
Description=Starting SUSI Server for SUSI Linux

[Service]
User=%i
# Tell SUSI server not to run as daemon, so that systemd can manage it better
Environment=DO_NOT_DAEMONIZE=1
ExecStart=@INSTALL_DIR@/bin/start.sh

[Install]
WantedBy=multi-user.target
