# systemd/ss-susi-server.service.in
[Unit]
Description=Starting SUSI Server for SUSI Linux

[Service]
User=@SUSI_SERVER_USER@
Type=forking
PIDFile=@INSTALL_DIR@/data/susi.pid
Environment=SS_LOG_TO_STDOUT=1
ExecStart=@INSTALL_DIR@/bin/start.sh

[Install]
WantedBy=multi-user.target
